# templates/default/mkswap.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        mkswap
# Required-Start:  
# Required-Stop:   
# Default-Start:   S
# Default-Stop:    
# Short-Description: format & start all swap attachements
### END INIT INFO

# *****************************************************************
# Dynamic iscsiadm configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
# Allows to take into account any resizing from the LVM cookbook...
PATH=/sbin:/bin:/usr/sbin:/usr/bin

swapList() { swapon -s| tail -n +2| cut -d' ' -f1; }

case "$1" in
start)	for i in $(swapList); do swapoff $i && mkswap -f $i && swapon $i; done ;;
*)
	echo "Syntaxe: $(basename $0) start" && false
esac
